# root/etc/services.d/unicorn-loadbalancer/run
exec s6-setuidgid abc /bin/sh -c 'npm start --prefix=/opt/UnicornLoadBalancer'